# laravel-worker.conf
# ================================================================
# Supervisor Configuration — Laravel Queue Worker
# النظام: إشعارات OneSignal الفورية عبر Redis
# ================================================================
#
# INSTALLATION على السيرفر (Ubuntu/Debian):
#   sudo apt-get install supervisor
#   sudo cp laravel-worker.conf /etc/supervisor/conf.d/laravel-worker.conf
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start laravel-worker:*
#
# ================================================================

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d

# ← غيّر المسار لمسار مشروعك على السيرفر
command=php /var/www/html/lab/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --queue=default

autostart=true
autorestart=true
stopasgroup=true
killasgroup=true

# عدد العمليات الموازية (زد لو عندك ضغط عالي)
numprocs=2

# ملفات اللوق
stdout_logfile=/var/www/html/lab/storage/logs/worker.log
stderr_logfile=/var/www/html/lab/storage/logs/worker-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

# المستخدم الذي يشغّل العملية (غيّره حسب السيرفر)
user=www-data

# إعادة تشغيل Worker بعد deploy جديد
stopwaitsecs=3600
